# HTTPServer.service
[Unit]
Description=HTTP Server

[Service]
Type=notify
ExecStart=/usr/bin/Debug/net6.0/HTTPServer.dll

[Install]
WantedBy=multi-user.target